# mongodb/supervisord-mongodb.conf
[program:mongodb]
command=mongod -f /data/mongodb.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=mongodb
autorestart=true
